# docs/pipewire-stereo-split/ytb-stereo-split.service
[Unit]
Description=Link OwnTone stereo-split loopback capture to the two AirPlay speakers
After=pipewire.service wireplumber.service
Requires=pipewire.service

[Service]
Type=simple
ExecStart=/usr/local/bin/ytb-stereo-split-linker.sh
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
